# LaunchAgents/co.nimblehq.cleanUpXcodeDerivedData.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>co.nimblehq.cleanUpXcodeDerivedData</string>
	<key>ProgramArguments</key>
	<array>
		<string>sh</string>
		<string>-c</string>
		<string>rm -rf ~/Library/Developer/Xcode/DerivedData &gt; ~/.log/clean_up_xcode_derived_data_logfile.txt 2&gt;&amp;1</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Day</key>
		<integer>2</integer>
	</dict>
</dict>
</plist>
